# twitter_bots/supervisor_config/conf.d/twitter_bots_prod.ini
[program:tweet-sender]
command=/usr/bin/python2.7 /home/robots/twitter_bots/manage.py tweet_sender --settings=twitter_bots.settings_prod
environment=PYTHONPATH=/home/robots/twitter_bots/
directory=/home/robots/twitter_bots/
user=robots
numprocs=1
stdout_logfile=/home/robots/twitter_bots/logs/supervisor/tweet_sender.log
stderr_logfile=/home/robots/twitter_bots/logs/supervisor/tweet_sender.errors.log
logfile_maxbytes=5MB
logfile_backups=5
autostart=true
autorestart=true
startsecs=1     ; número de segundos que el programa requiere estar funcionando para que se suponga que se arrancó bien
startretries=3     ; número máx de intentos para que el programa funcione bien 
stopwaitsecs = 10
priority=998


[program:create-tweets]
command=/usr/bin/python2.7 /home/robots/twitter_bots/manage.py create_tweets --settings=twitter_bots.settings_prod
environment=PYTHONPATH=/home/robots/twitter_bots/
directory=/home/robots/twitter_bots/
user=robots
numprocs=1
stdout_logfile=/home/robots/twitter_bots/logs/supervisor/create_tweets.log
stderr_logfile=/home/robots/twitter_bots/logs/supervisor/create_tweets.errors.log
logfile_maxbytes=5MB
logfile_backups=5
autostart=true
autorestart=true
startsecs=1     ; número de segundos que el programa requiere estar funcionando para que se suponga que se arrancó bien
startretries=3     ; número máx de intentos para que el programa funcione bien 
stopwaitsecs = 10
priority=998


[program:follower-extractor]
command=/usr/bin/python2.7 /home/robots/twitter_bots/manage.py follower_extractor --settings=twitter_bots.settings_prod
environment=PYTHONPATH=/home/robots/twitter_bots/
directory=/home/robots/twitter_bots/
user=robots
numprocs=1
stdout_logfile=/home/robots/twitter_bots/logs/supervisor/run_extractors.log
stderr_logfile=/home/robots/twitter_bots/logs/supervisor/run_extractors.errors.log
logfile_maxbytes=5MB
logfile_backups=5
autostart=true
autorestart=true
startsecs=1     ; número de segundos que el programa requiere estar funcionando para que se suponga que se arrancó bien
startretries=3     ; número máx de intentos para que el programa funcione bien 
stopwaitsecs = 10
priority=998


[program:bot-creation-finisher]
command=/usr/bin/python2.7 /home/robots/twitter_bots/manage.py bot_creation_finisher --settings=twitter_bots.settings_prod
environment=PYTHONPATH=/home/robots/twitter_bots/
directory=/home/robots/twitter_bots/
user=robots
numprocs=1
stdout_logfile=/home/robots/twitter_bots/logs/supervisor/bot_creation_finisher.log
stderr_logfile=/home/robots/twitter_bots/logs/supervisor/bot_creation_finisher.errors.log
logfile_maxbytes=5MB
logfile_backups=5
autostart=true
autorestart=true
startsecs=1     ; número de segundos que el programa requiere estar funcionando para que se suponga que se arrancó bien
startretries=3     ; número máx de intentos para que el programa funcione bien
stopwaitsecs = 10
priority=998
