# supervisor_config/conf.d/twitter_bots.ini
[program:follower-extractor]
command=/Users/Michel/RobotsTwitter/bin/python /Users/Michel/Documents/TwitterBots-REPO/manage.py follower_extractor
environment=PYTHONPATH=/Users/Michel/Documents/TwitterBots-REPO/
directory=/Users/Michel/Documents/TwitterBots-REPO/
numprocs=1
stdout_logfile=/Users/Michel/Documents/TwitterBots-REPO/logs/supervisor/follower_extractor.log
stderr_logfile=/Users/Michel/Documents/TwitterBots-REPO/logs/supervisor/follower_extractor.errors.log
logfile_maxbytes=5MB
logfile_backups=5
autostart=true
autorestart=true
startsecs=1     ; número de segundos que el programa requiere estar funcionando para que se suponga que se arrancó bien
startretries=3     ; número máx de intentos para que el programa funcione bien
stopwaitsecs = 10
priority=998


[program:bot-creator]
command=/Users/Michel/RobotsTwitter/bin/python /Users/Michel/Documents/TwitterBots-REPO/manage.py bot_creator
environment=PYTHONPATH=/Users/Michel/Documents/TwitterBots-REPO/
directory=/Users/Michel/Documents/TwitterBots-REPO/
numprocs=1
stdout_logfile=/Users/Michel/Documents/TwitterBots-REPO/logs/supervisor/bot_creator.log
stderr_logfile=/Users/Michel/Documents/TwitterBots-REPO/logs/supervisor/bot_creator.errors.log
logfile_maxbytes=5MB
logfile_backups=5
autostart=true
autorestart=true
startsecs=1     ; número de segundos que el programa requiere estar funcionando para que se suponga que se arrancó bien
startretries=3     ; número máx de intentos para que el programa funcione bien 
stopwaitsecs = 10
priority=998


[program:tweet-sender]
command=/Users/Michel/RobotsTwitter/bin/python /Users/Michel/Documents/TwitterBots-REPO/manage.py tweet_sender
environment=PYTHONPATH=/Users/Michel/Documents/TwitterBots-REPO/
directory=/Users/Michel/Documents/TwitterBots-REPO/
numprocs=1
stdout_logfile=/Users/Michel/Documents/TwitterBots-REPO/logs/supervisor/tweet_sender.log
stderr_logfile=/Users/Michel/Documents/TwitterBots-REPO/logs/supervisor/tweet_sender.errors.log
logfile_maxbytes=5MB
logfile_backups=5
autostart=true
autorestart=true
startsecs=1     ; número de segundos que el programa requiere estar funcionando para que se suponga que se arrancó bien
startretries=3     ; número máx de intentos para que el programa funcione bien 
stopwaitsecs = 10
priority=998